# heartbeat/config/com.heartbeat.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.heartbeat.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <!-- Update path below to your heartbeat directory -->
        <string>$HOME/.multiclaude/wts/orko/jolly-owl/heartbeat/run-heartbeat.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- Set your Anthropic API key here -->
        <key>ANTHROPIC_API_KEY</key>
        <string>YOUR_API_KEY_HERE</string>
        <key>HEARTBEAT_DIR</key>
        <string>$HOME/.multiclaude/wts/orko/jolly-owl/heartbeat</string>
    </dict>

    <!-- Run every 15 minutes (900 seconds) -->
    <key>StartInterval</key>
    <integer>900</integer>

    <key>StandardOutPath</key>
    <string>/tmp/heartbeat-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/heartbeat-stderr.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
